# autodeploy/cuty-rpc.conf
# /etc/init/cuty-rpc.conf

description "Web frontend to cutycapt"
author      "Tuomas Starck <tuomas@meetin.gs>"

start on    (local-filesystems and net-device-up IFACE!=lo)
stop  on    shutdown

respawn
respawn limit 5 5

setuid nobody
setgid nogroup

script
  set -e
  FIFO=`mktemp --dry-run`
  mkfifo $FIFO
  (logger -t cuty-rpc -p local1.notice <$FIFO &)
  exec 2>$FIFO
  exec >$FIFO
  rm $FIFO

  . /etc/autodeploy.conf

  export PORT=8000

  exec /usr/local/bin/$INTENT
end script
